# openwrt/etc/init.d/rvpn
#!/bin/sh /etc/rc.common
# Skvoz / rvpn

START=99
STOP=10
USE_PROCD=0

. /usr/lib/rvpn/common.sh
. /usr/lib/rvpn/dns.sh
. /usr/lib/rvpn/nft.sh
. /usr/lib/rvpn/singbox.sh
. /usr/lib/rvpn/zapret.sh
. /usr/lib/rvpn/health.sh
. /usr/lib/rvpn/watchdog.sh
. /usr/lib/rvpn/adblock.sh

# Fail-open cleanup shared by stop / hold / DNS failure
rvpn_runtime_failopen() {
	watchdog_stop
	nft_flush_zapret
	nft_flush_vpn
	nft_flush_quic
	nft_flush_doh
	ip rule del fwmark 0x1 lookup 100 2>/dev/null || true
	ip route flush table 100 2>/dev/null || true
	zapret_stop
	sb_kill_ours
	rm -f "$RVPN_RUN/sing-box.pid" "$RVPN_RUN/nfqws.pid" 2>/dev/null || true
	dns_heal_orphan 2>/dev/null || true
	dns_restore
	adblock_apply 2>/dev/null || true
}

start_service() {
	mkdir -p "$RVPN_RUN"
	chmod 700 "$RVPN_RUN" 2>/dev/null || true
	log "start requested"
	# UCI FakeIP survives reboot; /tmp backups do not — heal before anything else
	dns_heal_orphan 2>/dev/null || true
	ensure_ui_secret >/dev/null || log "WARN: ui_secret not ready"
	ensure_clash_secret >/dev/null || log "WARN: clash_secret not ready"
	health_cron_install 2>/dev/null || true

	# Soft failsafe hold: keep internet, do not auto-start VPN/zapret until cleared
	if rvpn_failsafe_hold_active && [ "${RVPN_CLEAR_HOLD:-0}" != "1" ]; then
		log "failsafe hold active — skip VPN/zapret (click Start / rvpnctl start to clear)"
		rvpn_runtime_failopen
		echo failsafe_hold >"$RVPN_STATE"
		return 0
	fi
	if [ "${RVPN_CLEAR_HOLD:-0}" = "1" ]; then
		rvpn_failsafe_hold_clear
		log "failsafe hold cleared"
	fi

	health_preflight || return 1

	zap=$(uci_get zapret_enabled)
	vpn=$(uci_get vpn_enabled)
	vpn_ok=0

	if [ "$zap" != "1" ] && [ "$vpn" != "1" ]; then
		log "both layers off — fail-open cleanup"
		rvpn_runtime_failopen
		echo stopped >"$RVPN_STATE"
		return 0
	fi

	if [ "$vpn" = "1" ]; then
		if sb_start; then
			nft_apply_vpn || log "WARN: nft vpn failed"
			vpn_ok=1
		else
			log "VPN start failed — fail-open (no FakeIP DNS)"
			nft_flush_vpn
		fi
	fi

	if [ "$zap" = "1" ]; then
		if zapret_start; then
			nft_apply_zapret || log "WARN: nft zapret failed"
		else
			log "zapret start failed"
		fi
	fi

	# FakeIP only if sing-box is up (dns_apply checks sb_pids)
	dns_apply
	nft_apply_quic
	nft_apply_doh
	adblock_apply 2>/dev/null || true

	if ! wan_ok; then
		log "CRITICAL: WAN broke — cleanup"
		stop_service
		return 1
	fi
	# DNS self-check: FakeIP must not be orphan after start
	if dns_uci_points_to_fakeip && ! dns_vpn_ready; then
		log "CRITICAL: FakeIP orphan after start — heal"
		dns_heal_orphan
	fi
	# Real resolve probe — catch "sb alive but DNS dead"
	if ! dns_lan_resolve_ok; then
		log "CRITICAL: DNS resolve failed after start — fail-open + hold"
		rvpn_runtime_failopen
		rvpn_failsafe_hold_set soft
		echo failsafe_hold >"$RVPN_STATE"
		return 1
	fi
	echo running >"$RVPN_STATE"
	if [ "$vpn_ok" = "1" ]; then
		watchdog_start
		if [ -f /usr/lib/rvpn/zapret-sync.sh ]; then
			if [ -f "$RVPN_RUN/zapret_sync.pending" ] || [ -f "$RVPN_RUN/nfqws_fetch.pending" ]; then
				# shellcheck source=/dev/null
				. /usr/lib/rvpn/zapret-sync.sh
				(
					sleep 5
					zapret_after_vpn_ready >>"$RVPN_LOG" 2>&1 || true
				) &
			fi
		fi
	else
		watchdog_stop
	fi
	log "start done vpn_ok=$vpn_ok"
}

stop_service() {
	log "stop requested"
	rvpn_runtime_failopen
	echo stopped >"$RVPN_STATE"
	log "stop done (fail-open)"
}

start() { start_service; }
stop() { stop_service; }
restart() {
	stop_service
	sleep 1
	# restart is always operator-initiated (boot uses start)
	RVPN_CLEAR_HOLD=1 start_service
}
