# openwrt/etc/init.d/rvpn
#!/bin/sh /etc/rc.common
# Skvoz / rvpn

START=99
STOP=10
USE_PROCD=0

. /usr/lib/rvpn/common.sh
. /usr/lib/rvpn/dns.sh
. /usr/lib/rvpn/nft.sh
. /usr/lib/rvpn/singbox.sh
. /usr/lib/rvpn/zapret.sh
. /usr/lib/rvpn/health.sh
. /usr/lib/rvpn/watchdog.sh

start_service() {
	mkdir -p "$RVPN_RUN"
	chmod 700 "$RVPN_RUN" 2>/dev/null || true
	log "start requested"
	ensure_ui_secret >/dev/null || log "WARN: ui_secret not ready"
	ensure_clash_secret >/dev/null || log "WARN: clash_secret not ready"
	health_preflight || return 1

	zap=$(uci_get zapret_enabled)
	vpn=$(uci_get vpn_enabled)
	vpn_ok=0

	if [ "$zap" != "1" ] && [ "$vpn" != "1" ]; then
		log "both layers off"
		echo stopped >"$RVPN_STATE"
		watchdog_stop
		return 0
	fi

	if [ "$vpn" = "1" ]; then
		if sb_start; then
			nft_apply_vpn || log "WARN: nft vpn failed"
			vpn_ok=1
		else
			log "VPN start failed — fail-open (no FakeIP DNS)"
			nft_flush_vpn
		fi
	fi

	if [ "$zap" = "1" ]; then
		if zapret_start; then
			nft_apply_zapret || log "WARN: nft zapret failed"
		else
			log "zapret start failed"
		fi
	fi

	# FakeIP only if sing-box is up (dns_apply checks sb_pids)
	dns_apply
	nft_apply_quic
	nft_apply_doh

	if ! wan_ok; then
		log "CRITICAL: WAN broke — cleanup"
		stop_service
		return 1
	fi
	echo running >"$RVPN_STATE"
	if [ "$vpn_ok" = "1" ]; then
		watchdog_start
	else
		watchdog_stop
	fi
	log "start done vpn_ok=$vpn_ok"
}

stop_service() {
	log "stop requested"
	watchdog_stop
	nft_flush_zapret
	nft_flush_vpn
	nft_flush_quic
	nft_flush_doh
	ip rule del fwmark 0x1 lookup 100 2>/dev/null || true
	ip route flush table 100 2>/dev/null || true
	zapret_stop
	sb_kill_ours
	rm -f "$RVPN_RUN/sing-box.pid" "$RVPN_RUN/nfqws.pid" 2>/dev/null || true
	dns_restore
	echo stopped >"$RVPN_STATE"
	log "stop done (fail-open)"
}

start() { start_service; }
stop() { stop_service; }
restart() {
	stop_service
	sleep 1
	start_service
}
